# qwen-api.service
[Unit]
Description=Qwen API
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=semo_db_user
Group=semo_db_user
WorkingDirectory=/mnt/storage/qwen-api
ExecStart=/mnt/storage/qwen-api/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1 --loop asyncio
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
KillMode=process

[Install]
WantedBy=multi-user.target